# tools/friday/plab/i3d.crontab
LOG_DIR=/home/ucb_i3/geels/
# every five minutes:
*/5 * * * *      cd $LOG_DIR && ./i3d.py clean >> $LOG_DIR/clean.log 2>&1
